# files/rc.local
/sbin/ifdown --force ens3
/sbin/ifup --force ens3
